# packer/webapp.service
Description=Webapp
After=network.target cloud-final.service
Requires=cloud-final.service

[Service]
ExecStart=/usr/bin/npm run start
WorkingDirectory=/home/admin/webapp
Restart=always
User=admin

[Install]
WantedBy=cloud-final.service